# deploy/wp-telegram-backup-bot.service.example
[Unit]
Description=Telegram WordPress Backup Bot
# Боту нужна исходящая сеть к Telegram API и локальная БД (MariaDB).
# Wants, а не Requires: падение/отсутствие БД не должно ронять бота целиком —
# он продолжит отвечать на /status, /disk; бэкап при недоступной БД просто завершится с ошибкой.
Wants=network-online.target mariadb.service
After=network-online.target mariadb.service
# Защита от crash-loop при битом .env/коде.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/wp-telegram-backup-bot
EnvironmentFile=/opt/wp-telegram-backup-bot/.env
ExecStart=/opt/wp-telegram-backup-bot/.venv/bin/python -m bot
Restart=on-failure
RestartSec=5
User=root
Group=root

# Бэкап-архивы и db.sql не должны быть доступны "прочим" пользователям.
UMask=0027

# Бэкап/restore могут идти долго — не убивать процесс мгновенно при stop.
TimeoutStopSec=120
KillMode=mixed

# Логи в journald (читать: journalctl -u wp-telegram-backup-bot -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wp-telegram-backup-bot

# Ограничения ресурсов:
# CPUQuota ограничивает компрессию zstd, чтобы не душить живой сайт во время бэкапа.
# Подстройте под сервер: 100% = одно ядро (для 1-ядерного VPS), на многоядерном
# можно дать 200% и больше. Уберите строку совсем, чтобы снять ограничение по CPU.
# Лимиты памяти и TasksMax намеренно НЕ задаём: корректный потолок памяти зависит от
# размера сайта и RAM сервера (слишком низкий OOM-убьёт бэкап), а боту нужно свободно
# порождать tar/mysqldump/runuser/chown. При необходимости добавьте сами.
CPUQuota=100%

[Install]
WantedBy=multi-user.target
